# agents/ceo/com.wse.ceo-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wse.ceo-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/npx</string>
    <string>tsx</string>
    <string>/Users/victoriaward/Projects/wse/agents/ceo/watcher.ts</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/victoriaward/Projects/wse/agents/ceo/logs/watcher.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/victoriaward/Projects/wse/agents/ceo/logs/watcher-error.log</string>
</dict>
</plist>
